# rootfs/etc/s6-overlay/s6-rc.d/addon-trmnl-screenshot/run
#!/command/execlineb -P
# Run the Node.js server
/usr/bin/node /app/index.js
